# logarchiver.service
[Unit]
Description=MeshCore Log Archiver Service - Archive and transfer logs weekly
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=YOUR_USERNAME
WorkingDirectory=/path/to/MeshMQTT
ExecStart=/path/to/MeshMQTT/venv/bin/python /path/to/MeshMQTT/archive_logs.py --month prev --send-each --no-compress --no-keep

[Install]
WantedBy=multi-user.target
