# mqttwatcher.service
[Unit]
Description=MeshCore MQTT Subscriber Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/MeshMQTT
ExecStart=/path/to/MeshMQTT/venv/bin/python /path/to/MeshMQTT/mqtt_subscriber.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
