# mqttwatcher.service
[Unit]
Description=MeshCore MQTT Subscriber Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
# User and WorkingDirectory will be replaced by install script
User=YOUR_USERNAME
WorkingDirectory=/path/to/meshcore/mqtt
ExecStart=/path/to/meshcore/mqtt/venv/bin/python /path/to/meshcore/mqtt/mqtt_subscriber.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/meshcore/mqtt/mqtt_logs/mqtt_subscriber.log
StandardError=append:/path/to/meshcore/mqtt/mqtt_logs/mqtt_subscriber_error.log

[Install]
WantedBy=multi-user.target
