# nodewatcher.service
[Unit]
Description=MeshCore Node Processor Service
After=network.target mqttwatcher.service
Requires=mqttwatcher.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/meshcore/mqtt
ExecStart=/path/to/meshcore/mqtt/venv/bin/python /path/to/meshcore/mqtt/meshnodes.py --watch --log mqtt_logs/data_log.jsonl --output nodes.json
Restart=always
RestartSec=10
[Install]
WantedBy=multi-user.target
